# linux-host/systemd/pasteover-server.service
[Unit]
Description=pasteover-server (Wayland clipboard over loopback for the builds VM)
Documentation=file:%h/pasteover/linux-host/HANDOFF.md
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/pasteover/linux-host/pasteover-server
Restart=always
RestartSec=2

[Install]
WantedBy=graphical-session.target
